# FocusDragonDaemon/Resources/com.focusdragon.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Unique identifier -->
    <key>Label</key>
    <string>com.focusdragon.daemon</string>

    <!-- Path to daemon executable -->
    <key>ProgramArguments</key>
    <array>
        <string>/Library/Application Support/FocusDragon/FocusDragonDaemon</string>
    </array>

    <!-- Run at system startup -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Keep daemon running, restart if crashes -->
    <key>KeepAlive</key>
    <dict>
        <!-- Always try to keep alive -->
        <key>SuccessfulExit</key>
        <false/>

        <!-- Restart after crashes -->
        <key>Crashed</key>
        <true/>
    </dict>

    <!-- Run as root user -->
    <key>UserName</key>
    <string>root</string>

    <!-- Run in root's group -->
    <key>GroupName</key>
    <string>wheel</string>

    <!-- Standard output (combined with daemon's own logging) -->
    <key>StandardOutPath</key>
    <string>/var/log/focusdragon/daemon-stdout.log</string>

    <!-- Error output -->
    <key>StandardErrorPath</key>
    <string>/var/log/focusdragon/daemon-stderr.log</string>

    <!-- Prevent rapid restart loops if daemon crashes repeatedly -->
    <key>ThrottleInterval</key>
    <integer>5</integer>

    <!-- Resource limits (prevent resource exhaustion) -->
    <key>HardResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>512</integer>
        <key>NumberOfProcesses</key>
        <integer>10</integer>
    </dict>

    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>256</integer>
    </dict>

    <!-- Lower priority (don't compete with critical system services) -->
    <key>Nice</key>
    <integer>10</integer>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Library/Application Support/FocusDragon</string>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
