# FocusDragonDaemon/Resources/com.focusdragon.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Label must match the plist filename (minus .plist) -->
    <key>Label</key>
    <string>com.focusdragon.daemon</string>

    <!-- SMAppService manages the actual path at runtime.
         BundleProgram tells it which executable inside the app bundle to use. -->
    <key>BundleProgram</key>
    <string>Contents/MacOS/FocusDragonDaemon</string>

    <!-- Run at system startup -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Keep daemon running, restart if crashes -->
    <key>KeepAlive</key>
    <true/>

    <!-- Standard output -->
    <key>StandardOutPath</key>
    <string>/var/log/focusdragon/daemon-stdout.log</string>

    <!-- Error output -->
    <key>StandardErrorPath</key>
    <string>/var/log/focusdragon/daemon-stderr.log</string>

    <!-- Prevent rapid restart loops -->
    <key>ThrottleInterval</key>
    <integer>5</integer>
</dict>
</plist>
